# puppet/modules/flume/files/etc/profile.d/flume-ng.sh
export FLUME_HOME=/usr/share/flume-ng
export PATH=$FLUME_HOME/bin:$PATH
