# modules/hadoop/files/etc/profile.d/hive.sh
export HIVE_HOME=/usr/share/hive
export PATH=$HIVE_HOME/bin:$PATH
